# debian/postinst
#!/bin/sh

# start soundcard-setup service
systemctl enable soundcard-setup.service

# make the software autostart
# TODO: this relies on the username being pi, and their home dir being in the standard place
if [ ! -L /home/pi/.config/autostart/pideck.desktop ] ; then
    mkdir -p /home/pi/.config/autostart
    chown pi:pi /home/pi/.config/autostart
    ln -s /usr/share/applications/pideck.desktop /home/pi/.config/autostart/pideck.desktop
fi

# reload the USB mounting udev rules
udevadm control --reload-rules

# remove the symlink added by usbmount
rm -rf /media/usb 2> /dev/null || :

# create a recordings directory and symlink it into /media
if [ ! -d /home/pi/recordings ] ; then
    mkdir -p /home/pi/recordings
    chown pi:pi /home/pi/recordings
    ln -s /home/pi/recordings /media/recordings
fi

# the script will always succeed
exit 0
